# deploy/systemd/openclaw-collector@.service
[Unit]
Description=OpenClaw Telegram collector (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fg
WorkingDirectory=/home/fg/openclaw_telegram_agent
EnvironmentFile=/etc/openclaw/%i.env
ExecStart=/home/fg/openclaw_telegram_agent/.venv/bin/python collector.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
